# systemd/system/getty@tty1.service.d/override.conf
# Autologin for user den
[Service]
ExecStart=
ExecStart=-/usr/bin/agetty --autologin den --noclear %I $TERM
